# RaspberryPi-FlightComputer/airplane.service
# /etc/systemd/system/airplane.service
[Unit]
Description=Damgaard Airplane Pi process
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=viktor
WorkingDirectory=/home/viktor/Projects/DamgaardAirplane2/RaspberryPi-FlightComputer

# Pre-start: update repo and reload units if the service file changed
ExecStartPre=/bin/bash -lc '\
  set -euo pipefail; \
  cd /home/viktor/Projects/DamgaardAirplane2/RaspberryPi-FlightComputer; \
  before="$(sha256sum airplane.service 2>/dev/null || true)"; \
  /usr/bin/git fetch origin main; \
  /usr/bin/git merge --ff-only origin/main || true; \
  after="$(sha256sum airplane.service 2>/dev/null || true)"; \
  if [ "$before" != "$after" ]; then /bin/systemctl daemon-reload; fi \
'

# Start your app
ExecStart=/usr/bin/env python3 /home/viktor/Projects/DamgaardAirplane2/RaspberryPi-FlightComputer/airplane_pi.py
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
